# ethereum/supervisor.d/geth.conf
[program:geth]
command=/home/ethereum/src/go-ethereum/build/bin/geth --datadir=/home/ethereum/.ethereum
autostart=true
user=ethereum
