# ethereum/supervisor.d/geth.conf
[program:geth]
command=/home/ethereum/src/go-ethereum/build/bin/geth --datadir=/home/ethereum/.ethereum --http --http.addr=0.0.0.0
autostart=true
user=ethereum
