# bitcoin/supervisor.d/bitcoind.conf
[program:bitcoind]
command=/usr/local/bitcoin/bin/bitcoind -printtoconsole -datadir=/home/bitcoin/.bitcoin -rpcbind=0.0.0.0 -rpcallowip=0.0.0.0/0
autostart=true
user=bitcoin
