# bitcoin/supervisor.d/bitcoind.conf
[program:bitcoind]
command=/usr/local/bitcoin/bin/bitcoind -printtoconsole -datadir=/home/bitcoin/.bitcoin
autostart=true
user=bitcoin
