# bitcoin/supervisor.d/ord-server.conf
[program:ord-server]
command=/usr/local/bin/ord server --address 0.0.0.0 --http --http-port 20546
environment=HOME=/home/bitcoin,USER=bitcoin
user=bitcoin
autostart=false
autorestart=false
