# bitcoin/supervisor.d/ord-index.conf
[program:ord-index]
command=/usr/local/bin/ord index
environment=HOME=/home/bitcoin,USER=bitcoin
user=bitcoin
autostart=true
autorestart=false
