# bitcoin/supervisor.d/ord.conf
[program:ord]
command=/usr/local/bin/ord server --address 127.0.0.1 --http --http-port 20546
environment=HOME=/home/bitcoin,USER=bitcoin
autostart=true
user=bitcoin
